# examples/timer/timer-example.service
[Unit]
Description=Inago Timer Example Service

[Service]
ExecStart=/bin/bash -c "echo Beep"
