# example/timer/timer-example.timer
[Unit]
Description=Inago Timer Example Timer

[Timer]
OnCalendar=weekly
Persistent=true     
 
[Install]
WantedBy=timers.target