# examples/mygroup/mygroup-bar@.service
[Unit]
Description=Inago Demo Unit

[Service]
ExecStart=/bin/bash -c "while true; do echo Hi; sleep 10; done"
